# ops/nanode-snapshots/nano-mumbai/openhydra-peer.service
[Unit]
Description=OpenHydra Peer (Qwen2.5-0.5B Shard 2, layers 12-16)
After=network.target
[Service]
Type=simple
WorkingDirectory=/opt/openhydra/src
Environment=MALLOC_ARENA_MAX=2
Environment=PYTHONPATH=/opt/openhydra/src
ExecStart=/opt/openhydra/.venv/bin/python3 -m coordinator.node --peer-id nano-mumbai --model-id openhydra-qwen2.5-0.5b --runtime-model-id /opt/openhydra/.cache/models/Qwen2.5-0.5B --runtime-backend pytorch_auto --layer-start 12 --layer-end 16 --shard-index 2 --total-shards 5 --grpc-port 50051 --api-port 8080 --api-host 0.0.0.0 --dht-url http://172.104.164.98:8468 --advertise-host 172.236.187.59
Restart=on-failure
RestartSec=10
MemoryMax=950M
LimitNOFILE=65536
[Install]
WantedBy=multi-user.target
